# ops/systemd/gbrain-profile.service
# gbrain · LinkedIn profile worker (Apify harvestapi scraper) — rich profiles,
# decoupled from the paid pipeline. Consumes gb_q_profile + backfills people who
# have a LinkedIn URL but no profile yet. No-ops without APIFY_TOKEN/APIFY_ACTOR_ID.
[Unit]
Description=gbrain LinkedIn profile worker (Apify harvestapi scraper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gbrain
WorkingDirectory=/opt/gbrain
EnvironmentFile=/opt/gbrain/.env
ExecStart=/opt/gbrain/.venv/bin/python -m workers.apify_linkedin worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
