# ops/systemd/gbrain-calendar.timer
# gbrain · run the Calendar connector every 5 minutes (syncToken → cheap polls)
# Install: copy both gbrain-calendar.{service,timer} to /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now gbrain-calendar.timer
# Check:   systemctl list-timers gbrain-calendar.timer
[Unit]
Description=Run gbrain Calendar connector every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
