# ops/systemd/gbrain-summarize.timer
# Weekly (opt-in: `systemctl enable --now gbrain-summarize.timer`).
[Unit]
Description=gbrain weekly org AI summaries
[Timer]
OnCalendar=Mon 09:00
Persistent=true
[Install]
WantedBy=timers.target
